# scripts/deploy/celery-worker.service
[Unit]
Description=Generative AI Vdos Celery worker (default queue)
After=network.target redis-server.service
Requires=redis-server.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/generative-ai-vdos
EnvironmentFile=/home/ubuntu/generative-ai-vdos/.env
Environment=DJANGO_SETTINGS_MODULE=config.settings.deployment
Environment=PYTHONPATH=backend
ExecStart=/home/ubuntu/generative-ai-vdos/.venv/bin/python -m celery -A config worker -l info -Q celery -n worker@%%h --concurrency 4
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
